# deploy/systemd/ignet-biosummarai.service
[Unit]
Description=Ignet BioSummarAI GPT-Powered Summarization Service
Documentation=https://github.com/hurlab/Ignet
After=network.target mariadb.service

[Service]
Type=simple
User=juhur
Group=webapps
WorkingDirectory=/data/var/www/html/ignet/biosummarAI
ExecStart=/data/var/www/html/ignet/biosummarAI/venv312/bin/python api_biosummary.py
EnvironmentFile=/data/var/www/html/ignet/biosummarAI/.env
Environment=PYTHONUNBUFFERED=1

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ignet-biosummarai

[Install]
WantedBy=multi-user.target
